# ffslfl-vpn-select/files/etc/hotplug.d/iface/60-vpn-select
#!/bin/sh
if [ "$ACTION" = "ifup" ] && [ "$INTERFACE" = "br-wan" ]
then
    sleep 3
    /usr/sbin/vpn-select
fi
