# ffslfl-vpn-select/files/etc/hotplug.d/iface/60-vpn-select
#!/bin/sh
[ "$ACTION" = ifup -a "$INTERFACE" = wan ] && {
       sleep 30
       /usr/sbin/vpn-select
}
