# ffslfl-vpn-select/files/etc/hotplug.d/iface/60-vpn-select
#!/bin/sh
[ "$ACTION" = "ifup" -a "$INTERFACE" = "br-wan" ] && {
       sleep 3
       /usr/sbin/vpn-select
}
